# login/95-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/home/USER/.config/scripts/hotplug"
